# scripts/DEBIAN/postinst
#!/bin/bash

systemctl daemon-reload
systemctl enable energieip-srv200-core.service
systemctl restart energieip-srv200-core.service
